# deploy/systemd/ihelped-api.service
[Unit]
Description=ihelped.ai API (Fastify + Node 20)
After=network.target mysql.service redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/ihelped.ai
EnvironmentFile=/etc/ihelped.env
ExecStart=/usr/bin/node /var/www/ihelped.ai/server/dist/index.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ReadWritePaths=/var/www/ihelped.ai
CapabilityBoundingSet=
LockPersonality=true
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
